# rootfs/etc/s6-overlay/s6-rc.d/omegabrr/run
#!/usr/bin/with-contenv sh

# run run (from /usr/local/bin)
exec omegabrr run
